# btsync.conf
[program:btsync]
user = nobody
command = /usr/sbin/btsync --nodaemon --config /config/btsync.conf